# systemd/.config/systemd/user/pull-projects.timer
[Unit]
Description=Pull repository updates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
