# systemd/.config/systemd/user/pull-repositories.service
[Unit]
Description=Pull repository updates daily

[Service]
ExecStart=%h/Repositories/pull %h/Repositories/
